# scripts/bin/services/suspend/suspend.service
[Unit]
Description=Run user script after suspend

[Service]
Type=oneshot
ExecStart=/home//bin/services/suspend/suspend.sh

[Install]
WantedBy=multi-user.target
