# dotlocal/bin/services/suspend/suspend.service
[Unit]
Description=Run user script after suspend

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/suspend.sh

[Install]
WantedBy=multi-user.target
